# ansible/roles/go2rtc/templates/go2rtc.service.j2
[Unit]
Description=go2rtc (RTSP -> WebRTC gateway)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/go2rtc -config /etc/go2rtc.yaml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
